# vk_to_tg.service
[Unit]
Description=VK to Telegram Bot (Docker Compose)
Documentation=https://github.com/j1mmy-n/vk-to-telegram-bot
Requires=docker.service
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes

# Измените путь, если репозиторий расположен не в /opt/vk-to-telegram-bot.
WorkingDirectory=/opt/vk-to-telegram-bot

ExecStart=/usr/bin/docker compose up -d --remove-orphans

# Reload загружает выбранную в .env версию образа и применяет её.
ExecReload=/usr/bin/docker compose pull
ExecReload=/usr/bin/docker compose up -d --remove-orphans

# Останавливает контейнеры, но сохраняет volumes с состоянием и логами.
ExecStop=/usr/bin/docker compose down

TimeoutStartSec=0
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
